# root/etc/services.d/overseerr/run
#!/usr/bin/with-contenv bash

cd /app/overseerr || exit

export CONFIG_DIRECTORY="/config"

exec \
    s6-setuidgid abc /usr/bin/yarn start
