# root/etc/s6-overlay/s6-rc.d/init-overseerr-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /run/overseerr-temp

if [[ -z ${LSIO_NON_ROOT_USER} ]]; then
    lsiown -R abc:abc \
        /config \
        /run/overseerr-temp
fi
